# blueprint/debian/debpkg/DEBIAN/prerm
#!/bin/bash -e
#
# Stop any daemons associated with the package before removal.
#

# remove mclient from systemd startup
systemctl disable mclient
