# debpkg/DEBIAN/prerm
#!/bin/bash -e

systemctl disable mclient
